# scroll-trim.service
[Unit]
Description=Trim scroll backups

[Service]
Type=oneshot
WorkingDirectory=/
ExecStart=/usr/bin/scroll trim . .
