# scroll-trim.timer
[Unit]
Description=Trigger scroll trimming

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
